# systemd/templates/node2-camera-control-agent.service.in
[Unit]
Description=Node2 AI Camera Control Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=@SERVICE_USER@
Group=@SERVICE_GROUP@
WorkingDirectory=@REPO_ROOT@
EnvironmentFile=@ENV_FILE@
ExecStartPre=@PYTHON@ scripts/common/render_runtime_config.py --role node2 --repo-root @REPO_ROOT@
ExecStartPre=@PYTHON@ scripts/common/validate_policy.py --policy @POLICY@
ExecStart=@UVICORN@ services.node2_control_agent.app:app --host @NODE2_IP@ --port @NODE2_API_PORT@ --no-proxy-headers
Restart=on-failure
RestartSec=3
UMask=0027
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
PrivateDevices=false
DeviceAllow=@DEVICE@ rw
ReadWritePaths=@RUNTIME_CONFIG_DIR@ @RESULTS_DIR@

[Install]
WantedBy=multi-user.target
