# systemd/templates/node1-ai-camera-receiver.service.in
[Unit]
Description=Node1 AI Camera RTP Receiver
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=@SERVICE_USER@
Group=@SERVICE_GROUP@
WorkingDirectory=@REPO_ROOT@
EnvironmentFile=@ENV_FILE@
ExecStartPre=@PYTHON@ scripts/common/apply_migrations.py --db-path @DB_PATH@ --migrations-dir @MIGRATIONS@
ExecStart=@PYTHON@ -m agents.node1.node1_receiver_agent --profile @PROFILE@ --transport @TRANSPORT@ --port @RTP_PORT@ --camera-id @CAMERA_ID@ --db-path @DB_PATH@ --metrics --metrics-port @NODE1_METRICS_PORT@ --motion-events --event-log @EVENT_LOG@ --latency-threshold-ms @LATENCY_THRESHOLD_MS@ --latency-window-samples @LATENCY_WINDOW_SAMPLES@ --no-exit-on-no-frames
Restart=on-failure
RestartSec=3
KillSignal=SIGINT
TimeoutStopSec=10
KillMode=mixed
FinalKillSignal=SIGKILL
UMask=0027
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=@DATA_DIR@ @RESULTS_DIR@

[Install]
WantedBy=multi-user.target
